# meta-renesas/meta-rz-distro/recipes-extend/v4l2-init/files/v4l2-init.service
[Unit]
Description=CRU initializing setting

[Service]
ExecStart=/root/v4l2-init.sh

[Install]
WantedBy=multi-user.target
